# templates/default/machine/init.d.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          <%= @init_script_name %>
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: AppDynamics machine agent
# Description:       AppDynamics machine agent
### END INIT INFO

# source function library
if [ -r /lib/lsb/init-functions ]; then
    . /lib/lsb/init-functions
elif [ -r /etc/rc.d/init.d/functions ]; then
    . /etc/rc.d/init.d/functions
elif [ -r /etc/init.d/functions ]; then
    . /etc/init.d/functions
fi

INSTALL_DIR="<%= @install_dir %>"
RETVAL=0
PROG="$INSTALL_DIR/run.sh"
pid_file="<%= @pid_file %>"

start() {
    echo -n $"Starting AppDynamics machine agent:"
    daemon --pidfile=${pid_file} ${PROG}
    RETVAL=$?
    echo
}

stop() {
    echo -n $"Stopping AppDynamics machine agent:"
    killproc -p ${pid_file} "$PROG" -TERM
    RETVAL=$?
    [ $? -eq 0 ] && rm -f ${pid_file}
    echo
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart}"
        RETVAL=1
        ;;
esac

exit $RETVAL
